# setup/collagen-tracking-worker.service
[Unit]
Description=Collagen Tracking Worker
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/collagen
Environment="COLLAGEN_DATA_DIR=/mnt/efs"
Environment="PYTHONUNBUFFERED=1"
ExecStart=/home/ubuntu/collagen/venv/bin/python3 /home/ubuntu/collagen/scripts/tracking_worker.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
